# sysNET/bsnet/units/targets/default.target
[Unit]
Description=Default target
After=multi-user.target