# sysNET/bsnet/units/targets/multi-user.target
[Unit]
Description=Multi User System
After=basic.target